# etc/init.d/fstrimDaemon
#!/sbin/runscript

checkconfig()
{
  if [ ! -f /etc/conf.d/fstrimDaemon ] ; then
    eerror "/etc/conf.d/fstrimDaemon does not exists."
    return 1
  fi

  if [ -z "${DAEMON}" ] \
  	|| [ -z "${PID}" ] \
  	|| [ -z "${IONICE}" ] \
  	|| [ -z "${NICE}" ] \
  	|| [ -z "${LOG}" ] \
  	|| [ -z "${TRIM_DIRS}" ] \
  	|| [ -z "${SLEEP_AT_START}" ] \
  	|| [ -z "${SLEEP_BEFORE_REPEAT}" ] \
  	|| [ -z "${SLEEP_BEFORE_REPEAT}" ]
  then
    eerror "You must set config options in /etc/conf.d/fstrimDaemon first"
    return 1
  fi

  if [ ! -f ${DAEMON} ] ; then
    eerror "Script ${DAEMON} does not exists."
    return 1
  fi
}

start()
{
  ebegin "Starting fstrimeDaemon"
  checkconfig || return 1
  echo `date`: Starting...  >> ${LOG}
  start-stop-daemon --start \
    --exec ${DAEMON} \
    --pidfile ${PID} \
    --background \
    --make-pidfile \
    --ionice ${IONICE} \
    --nice ${NICE} \
    --stdout ${LOG} \
    --stderr ${LOG}
  eend $?
}

stop()
{
  ebegin "Stopping fstrimeDaemon"
  start-stop-daemon --stop \
    --exec ${DAEMON} \
    --pidfile ${PID}
  eend $retval
  echo `date`: Stopped  >> ${LOG}
}
